# lumen-logs.timer
[Unit]
Description=Daily Lumen log creation timer

[Timer]
OnCalendar=*-*-* 02:00
Persistent=true

[Install]
WantedBy=timers.target
